# mailsystem/crontab
0 2 * * * postfix reload
0 2 * * * dovecot reload